# abcde.service
[Unit]
Description=ABCDE Rip

[Service]
ExecStart=/usr/local/bin/cd_rip.sh
ExecStopPost=/usr/bin/eject

[Install]
WantedBy=multi-user.target
